# MagicGatewayApp/deploy/magicgateway.service
[Unit]
Description=MagicGateway - DeepSeek API Gateway
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/beautycube/magicgateway
ExecStart=/home/beautycube/magicgateway/gateway
Restart=always
RestartSec=5
Environment=TZ=Asia/Shanghai

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
